# systemd/powalert.service
[Unit]
Description=System battery and power status notifier.

[Service]
Type=simple
ExecStart=/usr/bin/powalert

[Install]
WantedBy=graphical.target
